# upstart/my_blog.conf
description "my blog (ianww.com)"
start on filesystem and net-device-up IFACE=eth0
stop on shutdown

exec sudo -u ian "/home/ian/new_blog/prodstart.sh" >>/var/log/upstart/myblog.log
